# deploy/systemd/oink-digest.service
[Unit]
Description=Oink daily digest (search_jobs.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=oink
Group=oink
WorkingDirectory=/home/oink/app
EnvironmentFile=/home/oink/.env
Environment=PATH=/home/oink/venv/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/oink/app/skill/job-search/scripts
ExecStart=/home/oink/venv/bin/python skill/job-search/scripts/search_jobs.py
StandardOutput=journal
StandardError=journal

# 45 min — market_research workers + Opus calls can be slow when sources
# are degraded. If you tighten this, also tighten the per-source timeouts
# in skill/job-search/scripts/search_jobs.py.
TimeoutStartSec=2700

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/oink/state /home/oink/logs
ProtectKernelTunables=true
ProtectKernelModules=true
PrivateTmp=true
